# deploy/com.majicmovie.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.majicmovie</string>

    <!-- Adjust paths to your install location -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USER/Apps/MajicMovie</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USER/Apps/MajicMovie/.venv/bin/uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8443</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/majicmovie.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/majicmovie.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
